# private_dot_config/autostart/StartGpg.desktop
[Desktop Entry]
Encoding=UTF-8
Version=1
Type=Application
Name=StartGpg.desktop
Exec=/usr/bin/alacritty --title startGpg -e startGpg.sh
Hidden=false
